# pipeline/systemd/entertainment-ai-pipeline.service
[Unit]
Description=Entertainment AI daily batch pipeline
# Wait for real network (not just link-up) before starting, since we need
# to reach Postgres and Hugging Face.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/entertainment-ai/pipeline
EnvironmentFile=/home/ubuntu/entertainment-ai/pipeline/.env
ExecStart=/home/ubuntu/entertainment-ai/pipeline/run.sh

[Install]
WantedBy=multi-user.target
